# deploy/tencent-cloud/gangyantong-backend.service
[Unit]
Description=Gangyantong FastAPI backend
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/gangyantong/app/backend
EnvironmentFile=/opt/gangyantong/backend.env
ExecStart=/opt/gangyantong/app/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
